# debian/plug-ugmm.cron.daily
#!/bin/sh

php /usr/share/plug-ugmm/www/expiredmembers.php
