# debian/plug-ugmm.cron.daily
#!/bin/sh

php /usr/share/plug-ugmm/lib/expiredmembers.php
